# service-files/ck_health.timer
[Unit]
Description=Check ck-ospell health every 2 hours
[Timer]
OnCalendar=*-*-* 0/2:17:00
Persistent=true
[Install]
WantedBy=timers.target