# service-files/ck_restart.service
[Unit]
Description=Restart ck-ospell
[Service]
ExecStart=/bin/systemctl --user restart ck-ospell